# CPU_benchmarks/node-exporter/node-exporter.service
[Unit]
Description=Node Exporter Service
After=network.target
StartLimitIntervalSec=0

[Service]
LimitNOFILE=65536
Type=simple
Restart=always
RestartSec=1
ExecStart=/opt/node_exporter --collector.disable-defaults --collector.cpu --collector.cpufreq

[Install]
WantedBy=multi-user.target
